# github-follower.service
[Unit]
Description=GitHub Follower Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<your-username>
WorkingDirectory=/home/<your-username>/GithubFollower
EnvironmentFile=/home/<your-username>/GithubFollower/.env
ExecStart=/usr/bin/node /home/<your-username>/GithubFollower/node_modules/.bin/tsx src/daemon.ts
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-follower

[Install]
WantedBy=multi-user.target
